# scripts/systemd/auto_cert_require.timer
[Unit]
Description=Run cert auto acquire everyday

[Timer]
OnCalendar=04:19

[Install]
WantedBy=timers.target
