# scripts/systemd/auto-update-system.timer
[Unit]
Description=Run system-auto-update everyday

[Timer]
OnCalendar=02:19

[Install]
WantedBy=timers.target
