# scripts/systemd/auto-update-system.service
[Unit]
Description=Auto-update system by zypper

[Service]
ExecStart=/home/mio/iServ/scripts/auto_update.sh
