# scripts/systemd/auto_cert_require.service
[Unit]
Description=Auto acquire letsencrypt certs

[Service]
ExecStart=/home/mio/iServ/scripts/auto_certs_acquire.sh
