# scripts/systemd/nodequery.service
[Unit]
Description=Send ananysis data to nodequery

[Service]
ExecStart=/home/mio/iServ/scripts/invoke_nodequery.sh
