# scripts/systemd/wordpress_cronjob.service
[Unit]
Description=Invoke cronjob for wordpress

[Service]
ExecStart=/usr/bin/php -f /home/mio/www/blog/wordpress/wp-cron.php
