# scripts/systemd/wordpress_cronjob.timer
[Unit]
Description=Send ananysis data to nodequery every 180s

[Timer]
OnCalendar=30min

[Install]
WantedBy=timers.target
