# scripts/systemd/wordpress_cronjob.timer
[Unit]
Description=Invoke cronjob for wordpress ever hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
